# upstart/userjobs/templates/session_init_setup.conf
# {{ ansible_managed }}

description 'Setup User Sessions'

start on started dbus and runlevel [2345]
stop on runlevel [!2345]

task

script
  cat /etc/passwd | while read line; do
    user=`echo $line | cut -d: -f1`
    uid=`echo $line | cut -d: -f3`
    gid=`echo $line | cut -d: -f4`
    home=`echo $line | cut -d: -f6`

    if [ -d "$home/.init" ]; then
      # Create directory that would normally be
      # created by PAM when a user logs in.
      export XDG_RUNTIME_DIR="/run/user/$uid"
      mkdir -p "$XDG_RUNTIME_DIR"
      chmod 0700 "$XDG_RUNTIME_DIR"
      chown "$uid:$gid" "$XDG_RUNTIME_DIR"

      start session_init USER="$user"
    fi
  done
  initctl emit --no-wait user-jobs
end script
